# sigma_rules/offensive_tools/Telemetry/Telemetry.json
[
    {
        "keyword": "*TELEMETRY.exe install /url:*",
        "description": "Abusing Windows Telemetry for persistence through registry modifications and scheduled tasks to execute arbitrary commands with system-level privileges.",
        "tool_name": "Telemetry",
        "reference": "https://github.com/Imanfeng/Telemetry",
        "severity": "9",
        "popularity": "2",
        "keyword_type": "offensive_tool_keyword",
        "comment": "N/A",
        "tactics": "TA0003 - TA0005 - TA0004",
        "techniques": "T1053 - T1547 - T1059",
        "endpoint_detection": true,
        "network_detection": false
    },
    {
        "keyword": "*Execute command without file backdoor*",
        "description": "Abusing Windows Telemetry for persistence through registry modifications and scheduled tasks to execute arbitrary commands with system-level privileges.",
        "tool_name": "Telemetry",
        "reference": "https://github.com/Imanfeng/Telemetry",
        "severity": "9",
        "popularity": "2",
        "keyword_type": "offensive_tool_keyword",
        "comment": "N/A",
        "tactics": "TA0003 - TA0005 - TA0004",
        "techniques": "T1053 - T1547 - T1059",
        "endpoint_detection": true,
        "network_detection": false
    },
    {
        "keyword": "*5f026c27-f8e6-4052-b231-8451c6a73838*",
        "description": "Abusing Windows Telemetry for persistence through registry modifications and scheduled tasks to execute arbitrary commands with system-level privileges.",
        "tool_name": "Telemetry",
        "reference": "https://github.com/Imanfeng/Telemetry",
        "severity": "9",
        "popularity": "2",
        "keyword_type": "offensive_tool_keyword",
        "comment": "N/A",
        "tactics": "TA0003 - TA0005 - TA0004",
        "techniques": "T1053 - T1547 - T1059",
        "endpoint_detection": true,
        "network_detection": false
    },
    {
        "keyword": "* Remotely download Trojan files to *",
        "description": "Abusing Windows Telemetry for persistence through registry modifications and scheduled tasks to execute arbitrary commands with system-level privileges.",
        "tool_name": "Telemetry",
        "reference": "https://github.com/Imanfeng/Telemetry",
        "severity": "9",
        "popularity": "2",
        "keyword_type": "offensive_tool_keyword",
        "comment": "N/A",
        "tactics": "TA0003 - TA0005 - TA0004",
        "techniques": "T1053 - T1547 - T1059",
        "endpoint_detection": true,
        "network_detection": false
    },
    {
        "keyword": "*6455d210924926d364ebf88ff053821ff2d603ea99b17d1dbb454a7d061992cc*",
        "description": "Abusing Windows Telemetry for persistence through registry modifications and scheduled tasks to execute arbitrary commands with system-level privileges.",
        "tool_name": "Telemetry",
        "reference": "https://github.com/Imanfeng/Telemetry",
        "severity": "9",
        "popularity": "2",
        "keyword_type": "offensive_tool_keyword",
        "comment": "N/A",
        "tactics": "TA0003 - TA0005 - TA0004",
        "techniques": "T1053 - T1547 - T1059",
        "endpoint_detection": true,
        "network_detection": false
    },
    {
        "keyword": "*Hacktool.Msil.Telemetry*",
        "description": "Abusing Windows Telemetry for persistence through registry modifications and scheduled tasks to execute arbitrary commands with system-level privileges.",
        "tool_name": "Telemetry",
        "reference": "https://github.com/Imanfeng/Telemetry",
        "severity": "9",
        "popularity": "2",
        "keyword_type": "signature_keyword",
        "comment": "N/A",
        "tactics": "TA0003 - TA0005 - TA0004",
        "techniques": "T1053 - T1547 - T1059",
        "endpoint_detection": true,
        "network_detection": false
    },
    {
        "keyword": "*TELEMETRY.exe install /command:*",
        "description": "Abusing Windows Telemetry for persistence through registry modifications and scheduled tasks to execute arbitrary commands with system-level privileges.",
        "tool_name": "Telemetry",
        "reference": "https://github.com/Imanfeng/Telemetry",
        "severity": "9",
        "popularity": "2",
        "keyword_type": "offensive_tool_keyword",
        "comment": "N/A",
        "tactics": "TA0003 - TA0005 - TA0004",
        "techniques": "T1053 - T1547 - T1059",
        "endpoint_detection": true,
        "network_detection": false
    },
    {
        "keyword": "*\\SOFTWARE\\Microsoft\\Windows NT\\CurrentVersion\\AppCompatFlags\\TelemetryController\\fun*",
        "description": "Abusing Windows Telemetry for persistence through registry modifications and scheduled tasks to execute arbitrary commands with system-level privileges.",
        "tool_name": "Telemetry",
        "reference": "https://github.com/Imanfeng/Telemetry",
        "severity": "9",
        "popularity": "2",
        "keyword_type": "offensive_tool_keyword",
        "comment": "N/A",
        "tactics": "TA0003 - TA0005 - TA0004",
        "techniques": "T1053 - T1547 - T1059",
        "endpoint_detection": true,
        "network_detection": false
    },
    {
        "keyword": "*TELEMETRY.exe install /path:*",
        "description": "Abusing Windows Telemetry for persistence through registry modifications and scheduled tasks to execute arbitrary commands with system-level privileges.",
        "tool_name": "Telemetry",
        "reference": "https://github.com/Imanfeng/Telemetry",
        "severity": "9",
        "popularity": "2",
        "keyword_type": "offensive_tool_keyword",
        "comment": "N/A",
        "tactics": "TA0003 - TA0005 - TA0004",
        "techniques": "T1053 - T1547 - T1059",
        "endpoint_detection": true,
        "network_detection": false
    },
    {
        "keyword": "*https://www.trustedsec.com/blog/abusing-windows-telemetry-for-persistence/*",
        "description": "Abusing Windows Telemetry for persistence through registry modifications and scheduled tasks to execute arbitrary commands with system-level privileges.",
        "tool_name": "Telemetry",
        "reference": "https://github.com/Imanfeng/Telemetry",
        "severity": "9",
        "popularity": "2",
        "keyword_type": "offensive_tool_keyword",
        "comment": "N/A",
        "tactics": "TA0003 - TA0005 - TA0004",
        "techniques": "T1053 - T1547 - T1059",
        "endpoint_detection": true,
        "network_detection": true
    },
    {
        "keyword": "*2f00a05b-263d-4fcc-846b-da82bd684603*",
        "description": "Abusing Windows Telemetry for persistence through registry modifications and scheduled tasks to execute arbitrary commands with system-level privileges.",
        "tool_name": "Telemetry",
        "reference": "https://github.com/Imanfeng/Telemetry",
        "severity": "9",
        "popularity": "2",
        "keyword_type": "offensive_tool_keyword",
        "comment": "N/A",
        "tactics": "TA0003 - TA0005 - TA0004",
        "techniques": "T1053 - T1547 - T1059",
        "endpoint_detection": true,
        "network_detection": false
    },
    {
        "keyword": "*5F026C27-F8E6-4052-B231-8451C6A73838*",
        "description": "Abusing Windows Telemetry for persistence through registry modifications and scheduled tasks to execute arbitrary commands with system-level privileges.",
        "tool_name": "Telemetry",
        "reference": "https://github.com/Imanfeng/Telemetry",
        "severity": "9",
        "popularity": "2",
        "keyword_type": "offensive_tool_keyword",
        "comment": "N/A",
        "tactics": "TA0003 - TA0005 - TA0004",
        "techniques": "T1053 - T1547 - T1059",
        "endpoint_detection": true,
        "network_detection": false
    },
    {
        "keyword": "*58de3ab6935d1248e937e333e917586efb058e8b7d65ade38989543c806bd23e*",
        "description": "Abusing Windows Telemetry for persistence through registry modifications and scheduled tasks to execute arbitrary commands with system-level privileges.",
        "tool_name": "Telemetry",
        "reference": "https://github.com/Imanfeng/Telemetry",
        "severity": "9",
        "popularity": "2",
        "keyword_type": "offensive_tool_keyword",
        "comment": "N/A",
        "tactics": "TA0003 - TA0005 - TA0004",
        "techniques": "T1053 - T1547 - T1059",
        "endpoint_detection": true,
        "network_detection": false
    },
    {
        "keyword": "*/Telemetry.git*",
        "description": "Abusing Windows Telemetry for persistence through registry modifications and scheduled tasks to execute arbitrary commands with system-level privileges.",
        "tool_name": "Telemetry",
        "reference": "https://github.com/Imanfeng/Telemetry",
        "severity": "9",
        "popularity": "2",
        "keyword_type": "offensive_tool_keyword",
        "comment": "N/A",
        "tactics": "TA0003 - TA0005 - TA0004",
        "techniques": "T1053 - T1547 - T1059",
        "endpoint_detection": true,
        "network_detection": true
    },
    {
        "keyword": "*d7a308da069dcf3990f4cbfe57b8a1cc79c5f6b1259da795bba61592b8cf4b08*",
        "description": "Abusing Windows Telemetry for persistence through registry modifications and scheduled tasks to execute arbitrary commands with system-level privileges.",
        "tool_name": "Telemetry",
        "reference": "https://github.com/Imanfeng/Telemetry",
        "severity": "9",
        "popularity": "2",
        "keyword_type": "offensive_tool_keyword",
        "comment": "N/A",
        "tactics": "TA0003 - TA0005 - TA0004",
        "techniques": "T1053 - T1547 - T1059",
        "endpoint_detection": true,
        "network_detection": false
    },
    {
        "keyword": "*TELEMETRY.exe install *",
        "description": "Abusing Windows Telemetry for persistence through registry modifications and scheduled tasks to execute arbitrary commands with system-level privileges.",
        "tool_name": "Telemetry",
        "reference": "https://github.com/Imanfeng/Telemetry",
        "severity": "9",
        "popularity": "2",
        "keyword_type": "offensive_tool_keyword",
        "comment": "N/A",
        "tactics": "TA0003 - TA0005 - TA0004",
        "techniques": "T1053 - T1547 - T1059",
        "endpoint_detection": true,
        "network_detection": false
    },
    {
        "keyword": "*Imanfeng/Telemetry*",
        "description": "Abusing Windows Telemetry for persistence through registry modifications and scheduled tasks to execute arbitrary commands with system-level privileges.",
        "tool_name": "Telemetry",
        "reference": "https://github.com/Imanfeng/Telemetry",
        "severity": "9",
        "popularity": "2",
        "keyword_type": "offensive_tool_keyword",
        "comment": "N/A",
        "tactics": "TA0003 - TA0005 - TA0004",
        "techniques": "T1053 - T1547 - T1059",
        "endpoint_detection": true,
        "network_detection": true
    }
]