# sigma_rules/offensive_tools/SchTask_0x727/SchTask_0x727.json
[
    {
        "keyword": "*SchTask.exe C:\\Windows\\System32\\cmd.exe *",
        "description": "create hidden scheduled tasks",
        "tool_name": "SchTask_0x727",
        "reference": "https://github.com/0x727/SchTask_0x727",
        "severity": "10",
        "popularity": "6",
        "keyword_type": "offensive_tool_keyword",
        "comment": "N/A",
        "tactics": "TA0005 - TA0003",
        "techniques": "T1053",
        "endpoint_detection": true,
        "network_detection": false
    },
    {
        "keyword": "*0790530e1e0f1ed73b2b6fd701d75a2409c785af5367304d5fdd5bdfdf7eae46*",
        "description": "create hidden scheduled tasks",
        "tool_name": "SchTask_0x727",
        "reference": "https://github.com/0x727/SchTask_0x727",
        "severity": "10",
        "popularity": "6",
        "keyword_type": "offensive_tool_keyword",
        "comment": "#file_hash",
        "tactics": "TA0005 - TA0003",
        "techniques": "T1053",
        "endpoint_detection": true,
        "network_detection": false
    },
    {
        "keyword": "*24aae23bcf8b0a513988d69b1526eebd791007136a1faf08ea1df5a8d3884e50*",
        "description": "create hidden scheduled tasks",
        "tool_name": "SchTask_0x727",
        "reference": "https://github.com/0x727/SchTask_0x727",
        "severity": "10",
        "popularity": "6",
        "keyword_type": "offensive_tool_keyword",
        "comment": "#file_hash",
        "tactics": "TA0005 - TA0003",
        "techniques": "T1053",
        "endpoint_detection": true,
        "network_detection": false
    },
    {
        "keyword": "*b26b8713dc24bec3c5b0be456a1bbc058a8450c280d614695a691fa13ac6dbfd*",
        "description": "create hidden scheduled tasks",
        "tool_name": "SchTask_0x727",
        "reference": "https://github.com/0x727/SchTask_0x727",
        "severity": "10",
        "popularity": "6",
        "keyword_type": "offensive_tool_keyword",
        "comment": "#file_hash",
        "tactics": "TA0005 - TA0003",
        "techniques": "T1053",
        "endpoint_detection": true,
        "network_detection": false
    },
    {
        "keyword": "*/SchTask_0x727.git*",
        "description": "create hidden scheduled tasks",
        "tool_name": "SchTask_0x727",
        "reference": "https://github.com/0x727/SchTask_0x727",
        "severity": "10",
        "popularity": "6",
        "keyword_type": "offensive_tool_keyword",
        "comment": "N/A",
        "tactics": "TA0005 - TA0003",
        "techniques": "T1053",
        "endpoint_detection": true,
        "network_detection": true
    },
    {
        "keyword": "*] Hidden task xml file: *",
        "description": "create hidden scheduled tasks",
        "tool_name": "SchTask_0x727",
        "reference": "https://github.com/0x727/SchTask_0x727",
        "severity": "10",
        "popularity": "6",
        "keyword_type": "offensive_tool_keyword",
        "comment": "N/A",
        "tactics": "TA0005 - TA0003",
        "techniques": "T1053",
        "endpoint_detection": true,
        "network_detection": false
    },
    {
        "keyword": "*\\SchTask_0x727\\*",
        "description": "create hidden scheduled tasks",
        "tool_name": "SchTask_0x727",
        "reference": "https://github.com/0x727/SchTask_0x727",
        "severity": "10",
        "popularity": "6",
        "keyword_type": "offensive_tool_keyword",
        "comment": "N/A",
        "tactics": "TA0005 - TA0003",
        "techniques": "T1053",
        "endpoint_detection": true,
        "network_detection": false
    },
    {
        "keyword": "*/SchTask.zip*",
        "description": "create hidden scheduled tasks",
        "tool_name": "SchTask_0x727",
        "reference": "https://github.com/0x727/SchTask_0x727",
        "severity": "10",
        "popularity": "6",
        "keyword_type": "offensive_tool_keyword",
        "comment": "N/A",
        "tactics": "TA0005 - TA0003",
        "techniques": "T1053",
        "endpoint_detection": true,
        "network_detection": true
    },
    {
        "keyword": "*\\bin\\Release\\SchTask.exe*",
        "description": "create hidden scheduled tasks",
        "tool_name": "SchTask_0x727",
        "reference": "https://github.com/0x727/SchTask_0x727",
        "severity": "10",
        "popularity": "6",
        "keyword_type": "offensive_tool_keyword",
        "comment": "N/A",
        "tactics": "TA0005 - TA0003",
        "techniques": "T1053",
        "endpoint_detection": true,
        "network_detection": false
    },
    {
        "keyword": "*\\SchTask.zip*",
        "description": "create hidden scheduled tasks",
        "tool_name": "SchTask_0x727",
        "reference": "https://github.com/0x727/SchTask_0x727",
        "severity": "10",
        "popularity": "6",
        "keyword_type": "offensive_tool_keyword",
        "comment": "N/A",
        "tactics": "TA0005 - TA0003",
        "techniques": "T1053",
        "endpoint_detection": true,
        "network_detection": false
    },
    {
        "keyword": "*0x727/SchTask_0x727*",
        "description": "create hidden scheduled tasks",
        "tool_name": "SchTask_0x727",
        "reference": "https://github.com/0x727/SchTask_0x727",
        "severity": "10",
        "popularity": "6",
        "keyword_type": "offensive_tool_keyword",
        "comment": "N/A",
        "tactics": "TA0005 - TA0003",
        "techniques": "T1053",
        "endpoint_detection": true,
        "network_detection": true
    },
    {
        "keyword": "*70c104eb31780222a3a882a3728cafeeb308d8ff718a5c9ce62778d579b3de86*",
        "description": "create hidden scheduled tasks",
        "tool_name": "SchTask_0x727",
        "reference": "https://github.com/0x727/SchTask_0x727",
        "severity": "10",
        "popularity": "6",
        "keyword_type": "offensive_tool_keyword",
        "comment": "#file_hash",
        "tactics": "TA0005 - TA0003",
        "techniques": "T1053",
        "endpoint_detection": true,
        "network_detection": false
    },
    {
        "keyword": "*fe1ae959f9af863a11cefad541eb791a01e5bb9931cf5c57e478236ddad92ae6*",
        "description": "create hidden scheduled tasks",
        "tool_name": "SchTask_0x727",
        "reference": "https://github.com/0x727/SchTask_0x727",
        "severity": "10",
        "popularity": "6",
        "keyword_type": "offensive_tool_keyword",
        "comment": "#file_hash",
        "tactics": "TA0005 - TA0003",
        "techniques": "T1053",
        "endpoint_detection": true,
        "network_detection": false
    },
    {
        "keyword": "*SchTask_0x727/releases*",
        "description": "create hidden scheduled tasks",
        "tool_name": "SchTask_0x727",
        "reference": "https://github.com/0x727/SchTask_0x727",
        "severity": "10",
        "popularity": "6",
        "keyword_type": "offensive_tool_keyword",
        "comment": "N/A",
        "tactics": "TA0005 - TA0003",
        "techniques": "T1053",
        "endpoint_detection": true,
        "network_detection": true
    },
    {
        "keyword": "*E61C950E-A03D-40E2-AAD5-304C48570364*",
        "description": "create hidden scheduled tasks",
        "tool_name": "SchTask_0x727",
        "reference": "https://github.com/0x727/SchTask_0x727",
        "severity": "10",
        "popularity": "6",
        "keyword_type": "offensive_tool_keyword",
        "comment": "#GUID project",
        "tactics": "TA0005 - TA0003",
        "techniques": "T1053",
        "endpoint_detection": true,
        "network_detection": false
    },
    {
        "keyword": "*/SchTask_0x727/*",
        "description": "create hidden scheduled tasks",
        "tool_name": "SchTask_0x727",
        "reference": "https://github.com/0x727/SchTask_0x727",
        "severity": "10",
        "popularity": "6",
        "keyword_type": "offensive_tool_keyword",
        "comment": "N/A",
        "tactics": "TA0005 - TA0003",
        "techniques": "T1053",
        "endpoint_detection": true,
        "network_detection": true
    },
    {
        "keyword": "*\\SchTask.sln*",
        "description": "create hidden scheduled tasks",
        "tool_name": "SchTask_0x727",
        "reference": "https://github.com/0x727/SchTask_0x727",
        "severity": "10",
        "popularity": "6",
        "keyword_type": "offensive_tool_keyword",
        "comment": "N/A",
        "tactics": "TA0005 - TA0003",
        "techniques": "T1053",
        "endpoint_detection": true,
        "network_detection": false
    }
]