# sigma_rules/offensive_tools/Pspersist/Pspersist.json
[
    {
        "keyword": "*5A403F3C-9136-4B67-A94E-02D3BCD3162D*",
        "description": "Dropping a powershell script at %HOMEPATH%\\Documents\\windowspowershell\\ that contains the implant's path and whenever powershell process is created the implant will executed too.",
        "tool_name": "Pspersist",
        "reference": "https://github.com/TheD1rkMtr/Pspersist",
        "severity": "10",
        "popularity": "1",
        "keyword_type": "offensive_tool_keyword",
        "comment": "N/A",
        "tactics": "TA0003",
        "techniques": "T1546 - T1546.013 - T1053 - T1053.005 - T1037 - T1037.001",
        "endpoint_detection": true,
        "network_detection": false
    },
    {
        "keyword": "*Start Menu\\Programs\\Startup\\Loader.exe*",
        "description": "Dropping a powershell script at %HOMEPATH%\\Documents\\windowspowershell\\ that contains the implant's path and whenever powershell process is created the implant will executed too.",
        "tool_name": "Pspersist",
        "reference": "https://github.com/TheD1rkMtr/Pspersist",
        "severity": "10",
        "popularity": "1",
        "keyword_type": "offensive_tool_keyword",
        "comment": "N/A",
        "tactics": "TA0003",
        "techniques": "T1546 - T1546.013 - T1053 - T1053.005 - T1037 - T1037.001",
        "endpoint_detection": true,
        "network_detection": false
    },
    {
        "keyword": "*/PSpersist.git*",
        "description": "Dropping a powershell script at %HOMEPATH%\\Documents\\windowspowershell\\ that contains the implant's path and whenever powershell process is created the implant will executed too.",
        "tool_name": "Pspersist",
        "reference": "https://github.com/TheD1rkMtr/Pspersist",
        "severity": "10",
        "popularity": "1",
        "keyword_type": "offensive_tool_keyword",
        "comment": "N/A",
        "tactics": "TA0003",
        "techniques": "T1546 - T1546.013 - T1053 - T1053.005 - T1037 - T1037.001",
        "endpoint_detection": true,
        "network_detection": true
    },
    {
        "keyword": "*Pspersist-main*",
        "description": "Dropping a powershell script at %HOMEPATH%\\Documents\\windowspowershell\\ that contains the implant's path and whenever powershell process is created the implant will executed too.",
        "tool_name": "Pspersist",
        "reference": "https://github.com/TheD1rkMtr/Pspersist",
        "severity": "10",
        "popularity": "1",
        "keyword_type": "offensive_tool_keyword",
        "comment": "N/A",
        "tactics": "TA0003",
        "techniques": "T1546 - T1546.013 - T1053 - T1053.005 - T1037 - T1037.001",
        "endpoint_detection": true,
        "network_detection": true
    },
    {
        "keyword": "*TheD1rkMtr/Pspersist*",
        "description": "Dropping a powershell script at %HOMEPATH%\\Documents\\windowspowershell\\ that contains the implant's path and whenever powershell process is created the implant will executed too.",
        "tool_name": "Pspersist",
        "reference": "https://github.com/TheD1rkMtr/Pspersist",
        "severity": "10",
        "popularity": "1",
        "keyword_type": "offensive_tool_keyword",
        "comment": "N/A",
        "tactics": "TA0003",
        "techniques": "T1546 - T1546.013 - T1053 - T1053.005 - T1037 - T1037.001",
        "endpoint_detection": true,
        "network_detection": true
    },
    {
        "keyword": "*436b7f540f534a0ec1337cf82a76cb7727acda423132195f0c81560cdf75c438*",
        "description": "Dropping a powershell script at %HOMEPATH%\\Documents\\windowspowershell\\ that contains the implant's path and whenever powershell process is created the implant will executed too.",
        "tool_name": "Pspersist",
        "reference": "https://github.com/TheD1rkMtr/Pspersist",
        "severity": "10",
        "popularity": "1",
        "keyword_type": "offensive_tool_keyword",
        "comment": "N/A",
        "tactics": "TA0003",
        "techniques": "T1546 - T1546.013 - T1053 - T1053.005 - T1037 - T1037.001",
        "endpoint_detection": true,
        "network_detection": false
    },
    {
        "keyword": "*PSprofile.cpp*",
        "description": "Dropping a powershell script at %HOMEPATH%\\Documents\\windowspowershell\\ that contains the implant's path and whenever powershell process is created the implant will executed too.",
        "tool_name": "Pspersist",
        "reference": "https://github.com/TheD1rkMtr/Pspersist",
        "severity": "10",
        "popularity": "1",
        "keyword_type": "offensive_tool_keyword",
        "comment": "N/A",
        "tactics": "TA0003",
        "techniques": "T1546 - T1546.013 - T1053 - T1053.005 - T1037 - T1037.001",
        "endpoint_detection": true,
        "network_detection": false
    },
    {
        "keyword": "*58482e19d6376bbe0120289b6d39a35de15b68d00713f821ab0c7f28f85a31ee*",
        "description": "Dropping a powershell script at %HOMEPATH%\\Documents\\windowspowershell\\ that contains the implant's path and whenever powershell process is created the implant will executed too.",
        "tool_name": "Pspersist",
        "reference": "https://github.com/TheD1rkMtr/Pspersist",
        "severity": "10",
        "popularity": "1",
        "keyword_type": "offensive_tool_keyword",
        "comment": "N/A",
        "tactics": "TA0003",
        "techniques": "T1546 - T1546.013 - T1053 - T1053.005 - T1037 - T1037.001",
        "endpoint_detection": true,
        "network_detection": false
    },
    {
        "keyword": "*\\PSprofile.exe*",
        "description": "Dropping a powershell script at %HOMEPATH%\\Documents\\windowspowershell\\ that contains the implant's path and whenever powershell process is created the implant will executed too.",
        "tool_name": "Pspersist",
        "reference": "https://github.com/TheD1rkMtr/Pspersist",
        "severity": "10",
        "popularity": "1",
        "keyword_type": "offensive_tool_keyword",
        "comment": "N/A",
        "tactics": "TA0003",
        "techniques": "T1546 - T1546.013 - T1053 - T1053.005 - T1037 - T1037.001",
        "endpoint_detection": true,
        "network_detection": false
    }
]