# sigma_rules/greyware_tools/schtasks/schtasks.json
[
    {
        "keyword": "*schtasks /Create /RU SYSTEM /XML c:\\temp\\*",
        "description": "SSH backdoor creation with schtasks",
        "tool_name": "schtasks",
        "reference": "https://www.trellix.com/blogs/research/cactus-ransomware-new-strain-in-the-market/",
        "severity": "10",
        "popularity": "10",
        "keyword_type": "greyware_tool_keyword",
        "comment": "N/A",
        "tactics": "TA0003 - TA0005 - TA0011",
        "techniques": "T1053 - T1059.004 - T1090",
        "endpoint_detection": true,
        "network_detection": false
    },
    {
        "keyword": "* /Create /RU SYSTEM /TN MicrosoftEdgeUpdateTaskMachine /TR *",
        "description": "SSH backdoor creation with schtasks",
        "tool_name": "schtasks",
        "reference": "https://www.trellix.com/blogs/research/cactus-ransomware-new-strain-in-the-market/",
        "severity": "10",
        "popularity": "10",
        "keyword_type": "greyware_tool_keyword",
        "comment": "N/A",
        "tactics": "TA0003 - TA0005 - TA0011",
        "techniques": "T1053 - T1059.004 - T1090",
        "endpoint_detection": true,
        "network_detection": false
    },
    {
        "keyword": "*schtasks /query /v /fo LIST*",
        "description": "view detailed information about all the scheduled tasks.",
        "tool_name": "schtasks",
        "reference": "N/A",
        "severity": "8",
        "popularity": "10",
        "keyword_type": "greyware_tool_keyword",
        "comment": "N/A",
        "tactics": "TA0004 - TA0007",
        "techniques": "T1053.005 - T1082",
        "endpoint_detection": true,
        "network_detection": false
    },
    {
        "keyword": "*schtasks /Change /TN \"Microsoft\\Windows\\Windows Defender\\Windows Defender Verification\" /Disable*",
        "description": "disable scheduled tasks related to Windows Defender",
        "tool_name": "schtasks",
        "reference": "https://thedfirreport.com/2024/08/12/threat-actors-toolkit-leveraging-sliver-poshc2-batch-scripts/#c01",
        "severity": "10",
        "popularity": "10",
        "keyword_type": "greyware_tool_keyword",
        "comment": "N/A",
        "tactics": "TA0005",
        "techniques": "T1562.001",
        "endpoint_detection": true,
        "network_detection": false
    },
    {
        "keyword": "* /Create /TN sch.bat /TR \"c:\\temp\\script.vbs\" *",
        "description": "SSH backdoor creation with schtasks",
        "tool_name": "schtasks",
        "reference": "https://www.trellix.com/blogs/research/cactus-ransomware-new-strain-in-the-market/",
        "severity": "10",
        "popularity": "10",
        "keyword_type": "greyware_tool_keyword",
        "comment": "N/A",
        "tactics": "TA0003 - TA0005 - TA0011",
        "techniques": "T1053 - T1059.004 - T1090",
        "endpoint_detection": true,
        "network_detection": false
    },
    {
        "keyword": "*schtasks.exe /create /sc minute /mo 1 /tn * /rl highest /np /tr *\\sshd\\sshd.exe -f *\\sshd\\config\\sshd_config*",
        "description": "SSH backdoor creation with schtasks",
        "tool_name": "schtasks",
        "reference": "https://www.trellix.com/blogs/research/cactus-ransomware-new-strain-in-the-market/",
        "severity": "10",
        "popularity": "10",
        "keyword_type": "greyware_tool_keyword",
        "comment": "N/A",
        "tactics": "TA0003 - TA0005 - TA0011",
        "techniques": "T1053 - T1059.004 - T1090",
        "endpoint_detection": true,
        "network_detection": false
    },
    {
        "keyword": "*schtasks.exe /create /sc * /tr \"%programdata%\\sshd\\sshd.exe -f %programdata%\\sshd\\config\\sshd_config\\keys\\id_rsa -N -R * -o StrictHostKeyChecking=no -o *",
        "description": "SSH backdoor creation with schtasks",
        "tool_name": "schtasks",
        "reference": "https://www.trellix.com/blogs/research/cactus-ransomware-new-strain-in-the-market/",
        "severity": "10",
        "popularity": "10",
        "keyword_type": "greyware_tool_keyword",
        "comment": "N/A",
        "tactics": "TA0003 - TA0005 - TA0011",
        "techniques": "T1053 - T1059.004 - T1090",
        "endpoint_detection": true,
        "network_detection": false
    },
    {
        "keyword": "*schtasks /Change /TN \"Microsoft\\Windows\\Windows Defender\\Windows Defender Cache Maintenance\" /Disable*",
        "description": "disable scheduled tasks related to Windows Defender",
        "tool_name": "schtasks",
        "reference": "https://thedfirreport.com/2024/08/12/threat-actors-toolkit-leveraging-sliver-poshc2-batch-scripts/#c01",
        "severity": "10",
        "popularity": "10",
        "keyword_type": "greyware_tool_keyword",
        "comment": "N/A",
        "tactics": "TA0005",
        "techniques": "T1562.001",
        "endpoint_detection": true,
        "network_detection": false
    },
    {
        "keyword": "*schtasks /Change /TN \"Microsoft\\Windows\\Windows Defender\\Windows Defender Cleanup\" /Disable*",
        "description": "disable scheduled tasks related to Windows Defender",
        "tool_name": "schtasks",
        "reference": "https://thedfirreport.com/2024/08/12/threat-actors-toolkit-leveraging-sliver-poshc2-batch-scripts/#c01",
        "severity": "10",
        "popularity": "10",
        "keyword_type": "greyware_tool_keyword",
        "comment": "N/A",
        "tactics": "TA0005",
        "techniques": "T1562.001",
        "endpoint_detection": true,
        "network_detection": false
    },
    {
        "keyword": "*schtasks /Change /TN \"Microsoft\\Windows\\ExploitGuard\\ExploitGuard MDM policy Refresh\" /Disable*",
        "description": "disable scheduled tasks related to Windows Defender",
        "tool_name": "schtasks",
        "reference": "https://thedfirreport.com/2024/08/12/threat-actors-toolkit-leveraging-sliver-poshc2-batch-scripts/#c01",
        "severity": "10",
        "popularity": "10",
        "keyword_type": "greyware_tool_keyword",
        "comment": "N/A",
        "tactics": "TA0005",
        "techniques": "T1562.001",
        "endpoint_detection": true,
        "network_detection": false
    },
    {
        "keyword": "*schtasks /Change /TN \"Microsoft\\Windows\\Windows Defender\\Windows Defender Scheduled Scan\" /Disable*",
        "description": "disable scheduled tasks related to Windows Defender",
        "tool_name": "schtasks",
        "reference": "https://thedfirreport.com/2024/08/12/threat-actors-toolkit-leveraging-sliver-poshc2-batch-scripts/#c01",
        "severity": "10",
        "popularity": "10",
        "keyword_type": "greyware_tool_keyword",
        "comment": "N/A",
        "tactics": "TA0005",
        "techniques": "T1562.001",
        "endpoint_detection": true,
        "network_detection": false
    },
    {
        "keyword": "* /create /tn \"SysChecks\" /tr c:\\temp\\sch.bat *",
        "description": "SSH backdoor creation with schtasks",
        "tool_name": "schtasks",
        "reference": "https://www.trellix.com/blogs/research/cactus-ransomware-new-strain-in-the-market/",
        "severity": "10",
        "popularity": "10",
        "keyword_type": "greyware_tool_keyword",
        "comment": "N/A",
        "tactics": "TA0003 - TA0005 - TA0011",
        "techniques": "T1053 - T1059.004 - T1090",
        "endpoint_detection": true,
        "network_detection": false
    }
]