# sigma_rules/greyware_tools/schtasks/schtasks.json
[
    {
        "keyword": "*schtasks /query /v /fo LIST*",
        "description": "view detailed information about all the scheduled tasks.",
        "tool_name": "schtasks",
        "reference": "N/A",
        "severity": "8",
        "popularity": "10",
        "keyword_type": "greyware_tool_keyword",
        "comment": "N/A",
        "tactics": "TA0004 - TA0007",
        "techniques": "T1053.005 - T1082",
        "endpoint_detection": true,
        "network_detection": false
    }
]