# sigma_rules/greyware_tools/crontab/crontab.yml
title: Simple keyword detection rule for crontab
id: f753bb93-a3c3-4ce0-b520-814be8075f36
status: experimental
description: Detects interesting keywords based on crontab tool
references: []
author: '@mthcht'
date: 2023/07/30
modified: 2024/05/16
tags:
- attack.T1012
- attack.T1027
- attack.T1035
- attack.T1036
- attack.T1041
- attack.T1053.005
- attack.T1057
- attack.T1059.003
- attack.T1070.001
- attack.T1083
- attack.T1105
- attack.T1112
- attack.T1136
- attack.T1562.001
- attack.T1564.001
- attack.T1564.002
- attack.T1564.003
- attack.T1564.005
- attack.TA0002
- attack.TA0003
- attack.TA0006
- attack.TA0007
- attack.TA0008
- attack.TA0010
logsource:
  category:
  - endpoint
detection:
  selection:
  - ParentCommandLine:
    - '*crontab* sleep *ncat * -e /bin/bash*crontab*'
  - CommandLine:
    - '*crontab* sleep *ncat * -e /bin/bash*crontab*'
  - Details:
    - '*crontab* sleep *ncat * -e /bin/bash*crontab*'
  - registry.value:
    - '*crontab* sleep *ncat * -e /bin/bash*crontab*'
  - process.args:
    - '*crontab* sleep *ncat * -e /bin/bash*crontab*'
  - process.command_line:
    - '*crontab* sleep *ncat * -e /bin/bash*crontab*'
  - process.env_vars:
    - '*crontab* sleep *ncat * -e /bin/bash*crontab*'
  - process.io.text:
    - '*crontab* sleep *ncat * -e /bin/bash*crontab*'
  - process.title:
    - '*crontab* sleep *ncat * -e /bin/bash*crontab*'
  - pe.company:
    - '*crontab* sleep *ncat * -e /bin/bash*crontab*'
  - pe.description:
    - '*crontab* sleep *ncat * -e /bin/bash*crontab*'
  - pe.product:
    - '*crontab* sleep *ncat * -e /bin/bash*crontab*'
  - os.full:
    - '*crontab* sleep *ncat * -e /bin/bash*crontab*'
  - host.hostname:
    - '*crontab* sleep *ncat * -e /bin/bash*crontab*'
  - event.original:
    - '*crontab* sleep *ncat * -e /bin/bash*crontab*'
  - email.subject:
    - '*crontab* sleep *ncat * -e /bin/bash*crontab*'
  - device.model.name:
    - '*crontab* sleep *ncat * -e /bin/bash*crontab*'
  - container.image.name:
    - '*crontab* sleep *ncat * -e /bin/bash*crontab*'
  - container.name:
    - '*crontab* sleep *ncat * -e /bin/bash*crontab*'
  condition: selection
fields:
- ParentCommandLine
- CommandLine
- Details
- registry.value
- process.args
- process.command_line
- process.env_vars
- process.io.text
- process.title
- pe.company
- pe.description
- pe.product
- os.full
- host.hostname
- event.original
- email.subject
- device.model.name
- container.image.name
- container.name
falsepositives:
- unknown
level: medium
