# sigma_rules/offensive_tools/GhostTask/GhostTask.json
[
    {
        "keyword": "* - delete: Delete a scheduled task. Requires restarting the *",
        "description": "Creates scheduled tasks with a restrictive security descriptor -  making them invisible to all users. - Establishes scheduled tasks directly via the registry -  bypassing the generation of standard Windows event logs. - Provides support to modify existing scheduled tasks without generating Windows event logs. - Supports remote scheduled task creation (by using specially crafted Silver Ticket). - Supports to run in C2 with in-memory PE execution module (e.g. -  BruteRatel's memexec)",
        "tool_name": "GhostTask",
        "reference": "https://github.com/netero1010/GhostTask",
        "severity": "10",
        "popularity": "5",
        "keyword_type": "offensive_tool_keyword",
        "comment": "#contentstrings",
        "tactics": "TA0003 - TA0005 - TA0007",
        "techniques": "T1053.005 - T1112 - T1078",
        "endpoint_detection": true,
        "network_detection": false
    },
    {
        "keyword": "*a7ab668cab3a63df4a03cc53c46eed13fbb13bf1*",
        "description": "Creates scheduled tasks with a restrictive security descriptor -  making them invisible to all users. - Establishes scheduled tasks directly via the registry -  bypassing the generation of standard Windows event logs. - Provides support to modify existing scheduled tasks without generating Windows event logs. - Supports remote scheduled task creation (by using specially crafted Silver Ticket). - Supports to run in C2 with in-memory PE execution module (e.g. -  BruteRatel's memexec)",
        "tool_name": "GhostTask",
        "reference": "https://github.com/netero1010/GhostTask",
        "severity": "10",
        "popularity": "5",
        "keyword_type": "offensive_tool_keyword",
        "comment": "N/A",
        "tactics": "TA0003 - TA0005 - TA0007",
        "techniques": "T1053.005 - T1112 - T1078",
        "endpoint_detection": true,
        "network_detection": false
    },
    {
        "keyword": "*Successfully deleted scheduled task *",
        "description": "Creates scheduled tasks with a restrictive security descriptor -  making them invisible to all users. - Establishes scheduled tasks directly via the registry -  bypassing the generation of standard Windows event logs. - Provides support to modify existing scheduled tasks without generating Windows event logs. - Supports remote scheduled task creation (by using specially crafted Silver Ticket). - Supports to run in C2 with in-memory PE execution module (e.g. -  BruteRatel's memexec)",
        "tool_name": "GhostTask",
        "reference": "https://github.com/netero1010/GhostTask",
        "severity": "10",
        "popularity": "5",
        "keyword_type": "offensive_tool_keyword",
        "comment": "N/A",
        "tactics": "TA0003 - TA0005 - TA0007",
        "techniques": "T1053.005 - T1112 - T1078",
        "endpoint_detection": true,
        "network_detection": false
    },
    {
        "keyword": "* add * demon.x64.exe*",
        "description": "Creates scheduled tasks with a restrictive security descriptor -  making them invisible to all users. - Establishes scheduled tasks directly via the registry -  bypassing the generation of standard Windows event logs. - Provides support to modify existing scheduled tasks without generating Windows event logs. - Supports remote scheduled task creation (by using specially crafted Silver Ticket). - Supports to run in C2 with in-memory PE execution module (e.g. -  BruteRatel's memexec)",
        "tool_name": "GhostTask",
        "reference": "https://github.com/netero1010/GhostTask",
        "severity": "10",
        "popularity": "5",
        "keyword_type": "offensive_tool_keyword",
        "comment": "#contentstrings",
        "tactics": "TA0003 - TA0005 - TA0007",
        "techniques": "T1053.005 - T1112 - T1078",
        "endpoint_detection": true,
        "network_detection": false
    },
    {
        "keyword": "*.exe localhost add * \"cmd.exe\" \"/c *\" *logon*",
        "description": "Creates scheduled tasks with a restrictive security descriptor -  making them invisible to all users. - Establishes scheduled tasks directly via the registry -  bypassing the generation of standard Windows event logs. - Provides support to modify existing scheduled tasks without generating Windows event logs. - Supports remote scheduled task creation (by using specially crafted Silver Ticket). - Supports to run in C2 with in-memory PE execution module (e.g. -  BruteRatel's memexec)",
        "tool_name": "GhostTask",
        "reference": "https://github.com/netero1010/GhostTask",
        "severity": "10",
        "popularity": "5",
        "keyword_type": "offensive_tool_keyword",
        "comment": "N/A",
        "tactics": "TA0003 - TA0005 - TA0007",
        "techniques": "T1053.005 - T1112 - T1078",
        "endpoint_detection": true,
        "network_detection": false
    },
    {
        "keyword": "*\\GhostTask\\*",
        "description": "Creates scheduled tasks with a restrictive security descriptor -  making them invisible to all users. - Establishes scheduled tasks directly via the registry -  bypassing the generation of standard Windows event logs. - Provides support to modify existing scheduled tasks without generating Windows event logs. - Supports remote scheduled task creation (by using specially crafted Silver Ticket). - Supports to run in C2 with in-memory PE execution module (e.g. -  BruteRatel's memexec)",
        "tool_name": "GhostTask",
        "reference": "https://github.com/netero1010/GhostTask",
        "severity": "10",
        "popularity": "5",
        "keyword_type": "offensive_tool_keyword",
        "comment": "N/A",
        "tactics": "TA0003 - TA0005 - TA0007",
        "techniques": "T1053.005 - T1112 - T1078",
        "endpoint_detection": true,
        "network_detection": false
    },
    {
        "keyword": "*.exe localhost add * \"cmd.exe\" \"/c *\" *weekly*",
        "description": "Creates scheduled tasks with a restrictive security descriptor -  making them invisible to all users. - Establishes scheduled tasks directly via the registry -  bypassing the generation of standard Windows event logs. - Provides support to modify existing scheduled tasks without generating Windows event logs. - Supports remote scheduled task creation (by using specially crafted Silver Ticket). - Supports to run in C2 with in-memory PE execution module (e.g. -  BruteRatel's memexec)",
        "tool_name": "GhostTask",
        "reference": "https://github.com/netero1010/GhostTask",
        "severity": "10",
        "popularity": "5",
        "keyword_type": "offensive_tool_keyword",
        "comment": "N/A",
        "tactics": "TA0003 - TA0005 - TA0007",
        "techniques": "T1053.005 - T1112 - T1078",
        "endpoint_detection": true,
        "network_detection": false
    },
    {
        "keyword": "*DeleteKey(*SOFTWARE\\\\Microsoft\\\\Windows NT\\\\CurrentVersion\\\\Schedule\\\\TaskCache\\\\Tree\\\\*",
        "description": "Creates scheduled tasks with a restrictive security descriptor -  making them invisible to all users. - Establishes scheduled tasks directly via the registry -  bypassing the generation of standard Windows event logs. - Provides support to modify existing scheduled tasks without generating Windows event logs. - Supports remote scheduled task creation (by using specially crafted Silver Ticket). - Supports to run in C2 with in-memory PE execution module (e.g. -  BruteRatel's memexec)",
        "tool_name": "GhostTask",
        "reference": "https://github.com/netero1010/GhostTask",
        "severity": "10",
        "popularity": "5",
        "keyword_type": "offensive_tool_keyword",
        "comment": "N/A",
        "tactics": "TA0003 - TA0005 - TA0007",
        "techniques": "T1053.005 - T1112 - T1078",
        "endpoint_detection": true,
        "network_detection": false
    },
    {
        "keyword": "*\\GhostTask-1.0\\*",
        "description": "Creates scheduled tasks with a restrictive security descriptor -  making them invisible to all users. - Establishes scheduled tasks directly via the registry -  bypassing the generation of standard Windows event logs. - Provides support to modify existing scheduled tasks without generating Windows event logs. - Supports remote scheduled task creation (by using specially crafted Silver Ticket). - Supports to run in C2 with in-memory PE execution module (e.g. -  BruteRatel's memexec)",
        "tool_name": "GhostTask",
        "reference": "https://github.com/netero1010/GhostTask",
        "severity": "10",
        "popularity": "5",
        "keyword_type": "offensive_tool_keyword",
        "comment": "N/A",
        "tactics": "TA0003 - TA0005 - TA0007",
        "techniques": "T1053.005 - T1112 - T1078",
        "endpoint_detection": true,
        "network_detection": false
    },
    {
        "keyword": "* GhostTask.c *",
        "description": "Creates scheduled tasks with a restrictive security descriptor -  making them invisible to all users. - Establishes scheduled tasks directly via the registry -  bypassing the generation of standard Windows event logs. - Provides support to modify existing scheduled tasks without generating Windows event logs. - Supports remote scheduled task creation (by using specially crafted Silver Ticket). - Supports to run in C2 with in-memory PE execution module (e.g. -  BruteRatel's memexec)",
        "tool_name": "GhostTask",
        "reference": "https://github.com/netero1010/GhostTask",
        "severity": "10",
        "popularity": "5",
        "keyword_type": "offensive_tool_keyword",
        "comment": "N/A",
        "tactics": "TA0003 - TA0005 - TA0007",
        "techniques": "T1053.005 - T1112 - T1078",
        "endpoint_detection": true,
        "network_detection": false
    },
    {
        "keyword": "*netero1010/GhostTask*",
        "description": "Creates scheduled tasks with a restrictive security descriptor -  making them invisible to all users. - Establishes scheduled tasks directly via the registry -  bypassing the generation of standard Windows event logs. - Provides support to modify existing scheduled tasks without generating Windows event logs. - Supports remote scheduled task creation (by using specially crafted Silver Ticket). - Supports to run in C2 with in-memory PE execution module (e.g. -  BruteRatel's memexec)",
        "tool_name": "GhostTask",
        "reference": "https://github.com/netero1010/GhostTask",
        "severity": "10",
        "popularity": "5",
        "keyword_type": "offensive_tool_keyword",
        "comment": "N/A",
        "tactics": "TA0003 - TA0005 - TA0007",
        "techniques": "T1053.005 - T1112 - T1078",
        "endpoint_detection": true,
        "network_detection": true
    },
    {
        "keyword": "*DeleteScheduleTask(LPCSTR computerName*",
        "description": "Creates scheduled tasks with a restrictive security descriptor -  making them invisible to all users. - Establishes scheduled tasks directly via the registry -  bypassing the generation of standard Windows event logs. - Provides support to modify existing scheduled tasks without generating Windows event logs. - Supports remote scheduled task creation (by using specially crafted Silver Ticket). - Supports to run in C2 with in-memory PE execution module (e.g. -  BruteRatel's memexec)",
        "tool_name": "GhostTask",
        "reference": "https://github.com/netero1010/GhostTask",
        "severity": "10",
        "popularity": "5",
        "keyword_type": "offensive_tool_keyword",
        "comment": "N/A",
        "tactics": "TA0003 - TA0005 - TA0007",
        "techniques": "T1053.005 - T1112 - T1078",
        "endpoint_detection": true,
        "network_detection": false
    },
    {
        "keyword": "*GhostTask-1.0.zip*",
        "description": "Creates scheduled tasks with a restrictive security descriptor -  making them invisible to all users. - Establishes scheduled tasks directly via the registry -  bypassing the generation of standard Windows event logs. - Provides support to modify existing scheduled tasks without generating Windows event logs. - Supports remote scheduled task creation (by using specially crafted Silver Ticket). - Supports to run in C2 with in-memory PE execution module (e.g. -  BruteRatel's memexec)",
        "tool_name": "GhostTask",
        "reference": "https://github.com/netero1010/GhostTask",
        "severity": "10",
        "popularity": "5",
        "keyword_type": "offensive_tool_keyword",
        "comment": "N/A",
        "tactics": "TA0003 - TA0005 - TA0007",
        "techniques": "T1053.005 - T1112 - T1078",
        "endpoint_detection": true,
        "network_detection": false
    },
    {
        "keyword": "*.exe localhost add * \"cmd.exe\" \"/c *\" *second*",
        "description": "Creates scheduled tasks with a restrictive security descriptor -  making them invisible to all users. - Establishes scheduled tasks directly via the registry -  bypassing the generation of standard Windows event logs. - Provides support to modify existing scheduled tasks without generating Windows event logs. - Supports remote scheduled task creation (by using specially crafted Silver Ticket). - Supports to run in C2 with in-memory PE execution module (e.g. -  BruteRatel's memexec)",
        "tool_name": "GhostTask",
        "reference": "https://github.com/netero1010/GhostTask",
        "severity": "10",
        "popularity": "5",
        "keyword_type": "offensive_tool_keyword",
        "comment": "N/A",
        "tactics": "TA0003 - TA0005 - TA0007",
        "techniques": "T1053.005 - T1112 - T1078",
        "endpoint_detection": true,
        "network_detection": false
    },
    {
        "keyword": "*\\GhostTask-main*",
        "description": "Creates scheduled tasks with a restrictive security descriptor -  making them invisible to all users. - Establishes scheduled tasks directly via the registry -  bypassing the generation of standard Windows event logs. - Provides support to modify existing scheduled tasks without generating Windows event logs. - Supports remote scheduled task creation (by using specially crafted Silver Ticket). - Supports to run in C2 with in-memory PE execution module (e.g. -  BruteRatel's memexec)",
        "tool_name": "GhostTask",
        "reference": "https://github.com/netero1010/GhostTask",
        "severity": "10",
        "popularity": "5",
        "keyword_type": "offensive_tool_keyword",
        "comment": "N/A",
        "tactics": "TA0003 - TA0005 - TA0007",
        "techniques": "T1053.005 - T1112 - T1078",
        "endpoint_detection": true,
        "network_detection": false
    },
    {
        "keyword": "*GhostTask.exe*",
        "description": "Creates scheduled tasks with a restrictive security descriptor -  making them invisible to all users. - Establishes scheduled tasks directly via the registry -  bypassing the generation of standard Windows event logs. - Provides support to modify existing scheduled tasks without generating Windows event logs. - Supports remote scheduled task creation (by using specially crafted Silver Ticket). - Supports to run in C2 with in-memory PE execution module (e.g. -  BruteRatel's memexec)",
        "tool_name": "GhostTask",
        "reference": "https://github.com/netero1010/GhostTask",
        "severity": "10",
        "popularity": "5",
        "keyword_type": "offensive_tool_keyword",
        "comment": "N/A",
        "tactics": "TA0003 - TA0005 - TA0007",
        "techniques": "T1053.005 - T1112 - T1078",
        "endpoint_detection": true,
        "network_detection": true
    },
    {
        "keyword": "*.exe localhost add * \"cmd.exe\" \"/c *\" *daily*",
        "description": "Creates scheduled tasks with a restrictive security descriptor -  making them invisible to all users. - Establishes scheduled tasks directly via the registry -  bypassing the generation of standard Windows event logs. - Provides support to modify existing scheduled tasks without generating Windows event logs. - Supports remote scheduled task creation (by using specially crafted Silver Ticket). - Supports to run in C2 with in-memory PE execution module (e.g. -  BruteRatel's memexec)",
        "tool_name": "GhostTask",
        "reference": "https://github.com/netero1010/GhostTask",
        "severity": "10",
        "popularity": "5",
        "keyword_type": "offensive_tool_keyword",
        "comment": "N/A",
        "tactics": "TA0003 - TA0005 - TA0007",
        "techniques": "T1053.005 - T1112 - T1078",
        "endpoint_detection": true,
        "network_detection": false
    },
    {
        "keyword": "*/GhostTask.git*",
        "description": "Creates scheduled tasks with a restrictive security descriptor -  making them invisible to all users. - Establishes scheduled tasks directly via the registry -  bypassing the generation of standard Windows event logs. - Provides support to modify existing scheduled tasks without generating Windows event logs. - Supports remote scheduled task creation (by using specially crafted Silver Ticket). - Supports to run in C2 with in-memory PE execution module (e.g. -  BruteRatel's memexec)",
        "tool_name": "GhostTask",
        "reference": "https://github.com/netero1010/GhostTask",
        "severity": "10",
        "popularity": "5",
        "keyword_type": "offensive_tool_keyword",
        "comment": "N/A",
        "tactics": "TA0003 - TA0005 - TA0007",
        "techniques": "T1053.005 - T1112 - T1078",
        "endpoint_detection": true,
        "network_detection": true
    },
    {
        "keyword": "*\\GhostTask.h*",
        "description": "Creates scheduled tasks with a restrictive security descriptor -  making them invisible to all users. - Establishes scheduled tasks directly via the registry -  bypassing the generation of standard Windows event logs. - Provides support to modify existing scheduled tasks without generating Windows event logs. - Supports remote scheduled task creation (by using specially crafted Silver Ticket). - Supports to run in C2 with in-memory PE execution module (e.g. -  BruteRatel's memexec)",
        "tool_name": "GhostTask",
        "reference": "https://github.com/netero1010/GhostTask",
        "severity": "10",
        "popularity": "5",
        "keyword_type": "offensive_tool_keyword",
        "comment": "N/A",
        "tactics": "TA0003 - TA0005 - TA0007",
        "techniques": "T1053.005 - T1112 - T1078",
        "endpoint_detection": true,
        "network_detection": false
    }
]