# sigma_rules/offensive_tools/SharPersist/SharPersist.json
[
    {
        "keyword": "*INFO: Adding keepass backdoor persistence*",
        "description": "SharPersist Windows persistence toolkit written in C#.",
        "tool_name": "SharPersist",
        "reference": "https://github.com/fireeye/SharPersist",
        "severity": "10",
        "popularity": "10",
        "keyword_type": "offensive_tool_keyword",
        "comment": "N/A",
        "tactics": "TA0003 - TA0008",
        "techniques": "T1547 - T1053 - T1027 - T1028 - T1112",
        "endpoint_detection": true,
        "network_detection": false
    },
    {
        "keyword": "*.exe -t startupfolder -c * -a * -f*",
        "description": "SharPersist Windows persistence toolkit written in C#.",
        "tool_name": "SharPersist",
        "reference": "https://github.com/fireeye/SharPersist",
        "severity": "10",
        "popularity": "10",
        "keyword_type": "offensive_tool_keyword",
        "comment": "N/A",
        "tactics": "TA0003 - TA0008",
        "techniques": "T1547 - T1053 - T1027 - T1028 - T1112",
        "endpoint_detection": true,
        "network_detection": false
    },
    {
        "keyword": "*\\SharPersist\\*",
        "description": "SharPersist Windows persistence toolkit written in C#.",
        "tool_name": "SharPersist",
        "reference": "https://github.com/fireeye/SharPersist",
        "severity": "10",
        "popularity": "10",
        "keyword_type": "offensive_tool_keyword",
        "comment": "N/A",
        "tactics": "TA0003 - TA0008",
        "techniques": "T1547 - T1053 - T1027 - T1028 - T1112",
        "endpoint_detection": true,
        "network_detection": false
    },
    {
        "keyword": "*7806b81514ecc44219a6f6193b15b23aea0a947f3c91b339332bea1445745596*",
        "description": "SharPersist Windows persistence toolkit written in C#.",
        "tool_name": "SharPersist",
        "reference": "https://github.com/fireeye/SharPersist",
        "severity": "10",
        "popularity": "10",
        "keyword_type": "offensive_tool_keyword",
        "comment": "N/A",
        "tactics": "TA0003 - TA0008",
        "techniques": "T1547 - T1053 - T1027 - T1028 - T1112",
        "endpoint_detection": true,
        "network_detection": false
    },
    {
        "keyword": "*\\TortoiseSVNHookScripts.cs*",
        "description": "SharPersist Windows persistence toolkit written in C#.",
        "tool_name": "SharPersist",
        "reference": "https://github.com/fireeye/SharPersist",
        "severity": "10",
        "popularity": "10",
        "keyword_type": "offensive_tool_keyword",
        "comment": "N/A",
        "tactics": "TA0003 - TA0008",
        "techniques": "T1547 - T1053 - T1027 - T1028 - T1112",
        "endpoint_detection": true,
        "network_detection": false
    },
    {
        "keyword": "* -t schtaskbackdoor *",
        "description": "SharPersist Windows persistence toolkit written in C#.",
        "tool_name": "SharPersist",
        "reference": "https://github.com/fireeye/SharPersist",
        "severity": "10",
        "popularity": "10",
        "keyword_type": "offensive_tool_keyword",
        "comment": "N/A",
        "tactics": "TA0003 - TA0008",
        "techniques": "T1547 - T1053 - T1027 - T1028 - T1112",
        "endpoint_detection": true,
        "network_detection": false
    },
    {
        "keyword": "*INFO: Listing all scheduled tasks available to backdoor.*",
        "description": "SharPersist Windows persistence toolkit written in C#.",
        "tool_name": "SharPersist",
        "reference": "https://github.com/fireeye/SharPersist",
        "severity": "10",
        "popularity": "10",
        "keyword_type": "offensive_tool_keyword",
        "comment": "N/A",
        "tactics": "TA0003 - TA0008",
        "techniques": "T1547 - T1053 - T1027 - T1028 - T1112",
        "endpoint_detection": true,
        "network_detection": false
    },
    {
        "keyword": "*e9711f47cf9171f79bf34b342279f6fd9275c8ae65f3eb2c6ebb0b8432ea14f8*",
        "description": "SharPersist Windows persistence toolkit written in C#.",
        "tool_name": "SharPersist",
        "reference": "https://github.com/fireeye/SharPersist",
        "severity": "10",
        "popularity": "10",
        "keyword_type": "offensive_tool_keyword",
        "comment": "N/A",
        "tactics": "TA0003 - TA0008",
        "techniques": "T1547 - T1053 - T1027 - T1028 - T1112",
        "endpoint_detection": true,
        "network_detection": false
    },
    {
        "keyword": "*INFO: Checking backdoor present in KeePass config file*",
        "description": "SharPersist Windows persistence toolkit written in C#.",
        "tool_name": "SharPersist",
        "reference": "https://github.com/fireeye/SharPersist",
        "severity": "10",
        "popularity": "10",
        "keyword_type": "offensive_tool_keyword",
        "comment": "N/A",
        "tactics": "TA0003 - TA0008",
        "techniques": "T1547 - T1053 - T1027 - T1028 - T1112",
        "endpoint_detection": true,
        "network_detection": false
    },
    {
        "keyword": "*INFO: Adding service persistence*",
        "description": "SharPersist Windows persistence toolkit written in C#.",
        "tool_name": "SharPersist",
        "reference": "https://github.com/fireeye/SharPersist",
        "severity": "10",
        "popularity": "10",
        "keyword_type": "offensive_tool_keyword",
        "comment": "N/A",
        "tactics": "TA0003 - TA0008",
        "techniques": "T1547 - T1053 - T1027 - T1028 - T1112",
        "endpoint_detection": true,
        "network_detection": false
    },
    {
        "keyword": "*INFO: Adding scheduled task backdoor persistence*",
        "description": "SharPersist Windows persistence toolkit written in C#.",
        "tool_name": "SharPersist",
        "reference": "https://github.com/fireeye/SharPersist",
        "severity": "10",
        "popularity": "10",
        "keyword_type": "offensive_tool_keyword",
        "comment": "N/A",
        "tactics": "TA0003 - TA0008",
        "techniques": "T1547 - T1053 - T1027 - T1028 - T1112",
        "endpoint_detection": true,
        "network_detection": false
    },
    {
        "keyword": "*.exe -t tortoisesvn -c * -a * -m*",
        "description": "SharPersist Windows persistence toolkit written in C#.",
        "tool_name": "SharPersist",
        "reference": "https://github.com/fireeye/SharPersist",
        "severity": "10",
        "popularity": "10",
        "keyword_type": "offensive_tool_keyword",
        "comment": "N/A",
        "tactics": "TA0003 - TA0008",
        "techniques": "T1547 - T1053 - T1027 - T1028 - T1112",
        "endpoint_detection": true,
        "network_detection": false
    },
    {
        "keyword": "*SharPersist.exe*",
        "description": "SharPersist Windows persistence toolkit written in C#.",
        "tool_name": "SharPersist",
        "reference": "https://github.com/fireeye/SharPersist",
        "severity": "10",
        "popularity": "10",
        "keyword_type": "offensive_tool_keyword",
        "comment": "N/A",
        "tactics": "TA0003 - TA0008",
        "techniques": "T1547 - T1053 - T1027 - T1028 - T1112",
        "endpoint_detection": true,
        "network_detection": true
    },
    {
        "keyword": "*c:\\123.txt*",
        "description": "SharPersist Windows persistence toolkit written in C#.",
        "tool_name": "SharPersist",
        "reference": "https://github.com/fireeye/SharPersist",
        "severity": "10",
        "popularity": "10",
        "keyword_type": "offensive_tool_keyword",
        "comment": "N/A",
        "tactics": "TA0003 - TA0008",
        "techniques": "T1547 - T1053 - T1027 - T1028 - T1112",
        "endpoint_detection": true,
        "network_detection": false
    },
    {
        "keyword": "*INFO: Adding tortoise svn persistence*",
        "description": "SharPersist Windows persistence toolkit written in C#.",
        "tool_name": "SharPersist",
        "reference": "https://github.com/fireeye/SharPersist",
        "severity": "10",
        "popularity": "10",
        "keyword_type": "offensive_tool_keyword",
        "comment": "N/A",
        "tactics": "TA0003 - TA0008",
        "techniques": "T1547 - T1053 - T1027 - T1028 - T1112",
        "endpoint_detection": true,
        "network_detection": false
    },
    {
        "keyword": "*KeePassBackdoor.*",
        "description": "SharPersist Windows persistence toolkit written in C#.",
        "tool_name": "SharPersist",
        "reference": "https://github.com/fireeye/SharPersist",
        "severity": "10",
        "popularity": "10",
        "keyword_type": "offensive_tool_keyword",
        "comment": "N/A",
        "tactics": "TA0003 - TA0008",
        "techniques": "T1547 - T1053 - T1027 - T1028 - T1112",
        "endpoint_detection": true,
        "network_detection": false
    },
    {
        "keyword": "*9D1B853E-58F1-4BA5-AEFC-5C221CA30E48*",
        "description": "SharPersist Windows persistence toolkit written in C#.",
        "tool_name": "SharPersist",
        "reference": "https://github.com/fireeye/SharPersist",
        "severity": "10",
        "popularity": "10",
        "keyword_type": "offensive_tool_keyword",
        "comment": "N/A",
        "tactics": "TA0003 - TA0008",
        "techniques": "T1547 - T1053 - T1027 - T1028 - T1112",
        "endpoint_detection": true,
        "network_detection": false
    },
    {
        "keyword": "*.exe -t keepass -f *",
        "description": "SharPersist Windows persistence toolkit written in C#.",
        "tool_name": "SharPersist",
        "reference": "https://github.com/fireeye/SharPersist",
        "severity": "10",
        "popularity": "10",
        "keyword_type": "offensive_tool_keyword",
        "comment": "N/A",
        "tactics": "TA0003 - TA0008",
        "techniques": "T1547 - T1053 - T1027 - T1028 - T1112",
        "endpoint_detection": true,
        "network_detection": false
    },
    {
        "keyword": "*INFO: Adding registry persistence*",
        "description": "SharPersist Windows persistence toolkit written in C#.",
        "tool_name": "SharPersist",
        "reference": "https://github.com/fireeye/SharPersist",
        "severity": "10",
        "popularity": "10",
        "keyword_type": "offensive_tool_keyword",
        "comment": "N/A",
        "tactics": "TA0003 - TA0008",
        "techniques": "T1547 - T1053 - T1027 - T1028 - T1112",
        "endpoint_detection": true,
        "network_detection": false
    },
    {
        "keyword": "*SharPersist -*",
        "description": "SharPersist Windows persistence toolkit written in C#.",
        "tool_name": "SharPersist",
        "reference": "https://github.com/fireeye/SharPersist",
        "severity": "10",
        "popularity": "10",
        "keyword_type": "offensive_tool_keyword",
        "comment": "N/A",
        "tactics": "TA0003 - TA0008",
        "techniques": "T1547 - T1053 - T1027 - T1028 - T1112",
        "endpoint_detection": true,
        "network_detection": false
    },
    {
        "keyword": "* keepass backdoor persistence*",
        "description": "SharPersist Windows persistence toolkit written in C#.",
        "tool_name": "SharPersist",
        "reference": "https://github.com/fireeye/SharPersist",
        "severity": "10",
        "popularity": "10",
        "keyword_type": "offensive_tool_keyword",
        "comment": "N/A",
        "tactics": "TA0003 - TA0008",
        "techniques": "T1547 - T1053 - T1027 - T1028 - T1112",
        "endpoint_detection": true,
        "network_detection": false
    },
    {
        "keyword": "*INFO: Adding startup folder persistence*",
        "description": "SharPersist Windows persistence toolkit written in C#.",
        "tool_name": "SharPersist",
        "reference": "https://github.com/fireeye/SharPersist",
        "severity": "10",
        "popularity": "10",
        "keyword_type": "offensive_tool_keyword",
        "comment": "N/A",
        "tactics": "TA0003 - TA0008",
        "techniques": "T1547 - T1053 - T1027 - T1028 - T1112",
        "endpoint_detection": true,
        "network_detection": false
    },
    {
        "keyword": "* Keepass persistence backdoor *",
        "description": "SharPersist Windows persistence toolkit written in C#.",
        "tool_name": "SharPersist",
        "reference": "https://github.com/fireeye/SharPersist",
        "severity": "10",
        "popularity": "10",
        "keyword_type": "offensive_tool_keyword",
        "comment": "N/A",
        "tactics": "TA0003 - TA0008",
        "techniques": "T1547 - T1053 - T1027 - T1028 - T1112",
        "endpoint_detection": true,
        "network_detection": false
    },
    {
        "keyword": "*INFO: Adding scheduled task persistence*",
        "description": "SharPersist Windows persistence toolkit written in C#.",
        "tool_name": "SharPersist",
        "reference": "https://github.com/fireeye/SharPersist",
        "severity": "10",
        "popularity": "10",
        "keyword_type": "offensive_tool_keyword",
        "comment": "N/A",
        "tactics": "TA0003 - TA0008",
        "techniques": "T1547 - T1053 - T1027 - T1028 - T1112",
        "endpoint_detection": true,
        "network_detection": false
    },
    {
        "keyword": "*mandiant/SharPersist*",
        "description": "SharPersist Windows persistence toolkit written in C#.",
        "tool_name": "SharPersist",
        "reference": "https://github.com/fireeye/SharPersist",
        "severity": "10",
        "popularity": "10",
        "keyword_type": "offensive_tool_keyword",
        "comment": "N/A",
        "tactics": "TA0003 - TA0008",
        "techniques": "T1547 - T1053 - T1027 - T1028 - T1112",
        "endpoint_detection": true,
        "network_detection": true
    },
    {
        "keyword": "*SharPersist*",
        "description": "SharPersist Windows persistence toolkit written in C#.",
        "tool_name": "SharPersist",
        "reference": "https://github.com/fireeye/SharPersist",
        "severity": "10",
        "popularity": "10",
        "keyword_type": "offensive_tool_keyword",
        "comment": "N/A",
        "tactics": "TA0003 - TA0008",
        "techniques": "T1547 - T1053 - T1027 - T1028 - T1112",
        "endpoint_detection": true,
        "network_detection": false
    },
    {
        "keyword": "*\\SchTaskBackdoor.*",
        "description": "SharPersist Windows persistence toolkit written in C#.",
        "tool_name": "SharPersist",
        "reference": "https://github.com/fireeye/SharPersist",
        "severity": "10",
        "popularity": "10",
        "keyword_type": "offensive_tool_keyword",
        "comment": "N/A",
        "tactics": "TA0003 - TA0008",
        "techniques": "T1547 - T1053 - T1027 - T1028 - T1112",
        "endpoint_detection": true,
        "network_detection": false
    },
    {
        "keyword": "*KeePass config file is backdoored already*",
        "description": "SharPersist Windows persistence toolkit written in C#.",
        "tool_name": "SharPersist",
        "reference": "https://github.com/fireeye/SharPersist",
        "severity": "10",
        "popularity": "10",
        "keyword_type": "offensive_tool_keyword",
        "comment": "N/A",
        "tactics": "TA0003 - TA0008",
        "techniques": "T1547 - T1053 - T1027 - T1028 - T1112",
        "endpoint_detection": true,
        "network_detection": false
    },
    {
        "keyword": "*/SharPersist.git*",
        "description": "SharPersist Windows persistence toolkit written in C#.",
        "tool_name": "SharPersist",
        "reference": "https://github.com/fireeye/SharPersist",
        "severity": "10",
        "popularity": "10",
        "keyword_type": "offensive_tool_keyword",
        "comment": "N/A",
        "tactics": "TA0003 - TA0008",
        "techniques": "T1547 - T1053 - T1027 - T1028 - T1112",
        "endpoint_detection": true,
        "network_detection": true
    },
    {
        "keyword": "*SchTaskBackdoor.*",
        "description": "SharPersist Windows persistence toolkit written in C#.",
        "tool_name": "SharPersist",
        "reference": "https://github.com/fireeye/SharPersist",
        "severity": "10",
        "popularity": "10",
        "keyword_type": "offensive_tool_keyword",
        "comment": "N/A",
        "tactics": "TA0003 - TA0008",
        "techniques": "T1547 - T1053 - T1027 - T1028 - T1112",
        "endpoint_detection": true,
        "network_detection": false
    }
]