# sigma_rules/greyware_tools/nc/nc.json
[
    {
        "keyword": "*nc -l -p * -e /bin/bash*",
        "description": "Netcat Backdoor on Linux - create a relay that sends packets from the local port to a netcat client connecte to the target ip on the targeted port",
        "tool_name": "nc",
        "reference": "https://github.com/RoseSecurity/Red-Teaming-TTPs/blob/main/NetcatCheatSheet.pdf",
        "severity": "N/A",
        "popularity": "10",
        "keyword_type": "greyware_tool_keyword",
        "comment": "N/A",
        "tactics": "TA0011 - TA0040",
        "techniques": "T1090.001 - T1021.001",
        "endpoint_detection": true,
        "network_detection": false
    },
    {
        "keyword": "*nc -z -v * *",
        "description": "netcat common arguments",
        "tool_name": "nc",
        "reference": "N/A",
        "severity": "N/A",
        "popularity": "N/A",
        "keyword_type": "greyware_tool_keyword",
        "comment": "N/A",
        "tactics": "TA0011 - TA0040",
        "techniques": "T1090.001 - T1021.001",
        "endpoint_detection": true,
        "network_detection": false
    },
    {
        "keyword": "* /bin/nc * -e /bin/bash* > cron && crontab cron*",
        "description": "Linux Persistence Shell cron",
        "tool_name": "nc",
        "reference": "https://github.com/RoseSecurity/Red-Teaming-TTPs/blob/main/Linux.md",
        "severity": "N/A",
        "popularity": "10",
        "keyword_type": "greyware_tool_keyword",
        "comment": "N/A",
        "tactics": "TA0003",
        "techniques": "T1053 - T1037",
        "endpoint_detection": true,
        "network_detection": false
    },
    {
        "keyword": "*nc -l -p * -e cmd.exe*",
        "description": "Netcat Backdoor on Windows - create a relay that sends packets from the local port to a netcat client connecte to the target ip on the targeted port",
        "tool_name": "nc",
        "reference": "https://github.com/RoseSecurity/Red-Teaming-TTPs/blob/main/NetcatCheatSheet.pdf",
        "severity": "N/A",
        "popularity": "10",
        "keyword_type": "greyware_tool_keyword",
        "comment": "N/A",
        "tactics": "TA0011 - TA0040",
        "techniques": "T1090.001 - T1021.001",
        "endpoint_detection": true,
        "network_detection": false
    },
    {
        "keyword": "*echo nc -l -p * > *.bat*",
        "description": "Netcat Realy on windows - create a relay that sends packets from the local port to a netcat client connecte to the target ip on the targeted port",
        "tool_name": "nc",
        "reference": "https://github.com/RoseSecurity/Red-Teaming-TTPs/blob/main/NetcatCheatSheet.pdf",
        "severity": "N/A",
        "popularity": "10",
        "keyword_type": "greyware_tool_keyword",
        "comment": "N/A",
        "tactics": "TA0011 - TA0040",
        "techniques": "T1090.001 - T1021.001",
        "endpoint_detection": true,
        "network_detection": false
    },
    {
        "keyword": "*nc -v -n -z -w1 *-*",
        "description": "Port scanner with netcat",
        "tool_name": "nc",
        "reference": "https://github.com/RoseSecurity/Red-Teaming-TTPs/blob/main/NetcatCheatSheet.pdf",
        "severity": "N/A",
        "popularity": "10",
        "keyword_type": "greyware_tool_keyword",
        "comment": "N/A",
        "tactics": "TA0007",
        "techniques": "T1046",
        "endpoint_detection": true,
        "network_detection": false
    },
    {
        "keyword": "* /bin/nc * -e /bin/bash*> * crontab cron*",
        "description": "linux commands abused by attackers",
        "tool_name": "nc",
        "reference": "N/A",
        "severity": "N/A",
        "popularity": "N/A",
        "keyword_type": "greyware_tool_keyword",
        "comment": "greyware_tools high risks of false positives",
        "tactics": "TA0003 - TA0007 - TA0008 - TA0010 - TA0006 - TA0002",
        "techniques": "T1059.003 - T1053.005 - T1105 - T1012 - T1057 - T1083 - T1041 - T1036 - T1035 - T1562.001 - T1564.001 - T1564.005 - T1564.002 - T1564.003 - T1027 - T1070.001 - T1112 - T1136",
        "endpoint_detection": true,
        "network_detection": false
    },
    {
        "keyword": "*nc -l -p * -e *.bat*",
        "description": "Netcat Realy on windows - create a relay that sends packets from the local port to a netcat client connecte to the target ip on the targeted port",
        "tool_name": "nc",
        "reference": "https://github.com/RoseSecurity/Red-Teaming-TTPs/blob/main/NetcatCheatSheet.pdf",
        "severity": "N/A",
        "popularity": "10",
        "keyword_type": "greyware_tool_keyword",
        "comment": "N/A",
        "tactics": "TA0011 - TA0040",
        "techniques": "T1090.001 - T1021.001",
        "endpoint_detection": true,
        "network_detection": false
    }
]