# sigma_rules/greyware_tools/nc/nc.yml
title: Simple keyword detection rule for nc
id: 43122976-38a1-4bd5-afd7-3730a2fce096
status: experimental
description: Detects interesting keywords based on nc tool
references:
- https://github.com/RoseSecurity/Red-Teaming-TTPs/blob/main/Linux.md
- https://github.com/RoseSecurity/Red-Teaming-TTPs/blob/main/NetcatCheatSheet.pdf
author: '@mthcht'
date: 2023/07/30
modified: 2024/07/01
tags:
- attack.T1012
- attack.T1021.001
- attack.T1027
- attack.T1035
- attack.T1036
- attack.T1037
- attack.T1041
- attack.T1046
- attack.T1053
- attack.T1053.005
- attack.T1057
- attack.T1059.003
- attack.T1070.001
- attack.T1083
- attack.T1090.001
- attack.T1105
- attack.T1112
- attack.T1136
- attack.T1562.001
- attack.T1564.001
- attack.T1564.002
- attack.T1564.003
- attack.T1564.005
- attack.TA0002
- attack.TA0003
- attack.TA0006
- attack.TA0007
- attack.TA0008
- attack.TA0010
- attack.TA0011
- attack.TA0040
logsource:
  category:
  - endpoint
detection:
  selection:
  - ParentCommandLine:
    - '*nc -l -p * -e cmd.exe*'
    - '*nc -l -p * -e *.bat*'
    - '*nc -l -p * -e /bin/bash*'
    - '*echo nc -l -p * > *.bat*'
    - '*nc -v -n -z -w1 *-*'
    - '*nc -z -v * *'
    - '* /bin/nc * -e /bin/bash* > cron && crontab cron*'
    - '* /bin/nc * -e /bin/bash*> * crontab cron*'
  - CommandLine:
    - '*nc -l -p * -e cmd.exe*'
    - '*nc -l -p * -e *.bat*'
    - '*nc -l -p * -e /bin/bash*'
    - '*echo nc -l -p * > *.bat*'
    - '*nc -v -n -z -w1 *-*'
    - '*nc -z -v * *'
    - '* /bin/nc * -e /bin/bash* > cron && crontab cron*'
    - '* /bin/nc * -e /bin/bash*> * crontab cron*'
  - Details:
    - '*nc -l -p * -e cmd.exe*'
    - '*nc -l -p * -e *.bat*'
    - '*nc -l -p * -e /bin/bash*'
    - '*echo nc -l -p * > *.bat*'
    - '*nc -v -n -z -w1 *-*'
    - '*nc -z -v * *'
    - '* /bin/nc * -e /bin/bash* > cron && crontab cron*'
    - '* /bin/nc * -e /bin/bash*> * crontab cron*'
  - registry.value:
    - '*nc -l -p * -e cmd.exe*'
    - '*nc -l -p * -e *.bat*'
    - '*nc -l -p * -e /bin/bash*'
    - '*echo nc -l -p * > *.bat*'
    - '*nc -v -n -z -w1 *-*'
    - '*nc -z -v * *'
    - '* /bin/nc * -e /bin/bash* > cron && crontab cron*'
    - '* /bin/nc * -e /bin/bash*> * crontab cron*'
  - process.args:
    - '*nc -l -p * -e cmd.exe*'
    - '*nc -l -p * -e *.bat*'
    - '*nc -l -p * -e /bin/bash*'
    - '*echo nc -l -p * > *.bat*'
    - '*nc -v -n -z -w1 *-*'
    - '*nc -z -v * *'
    - '* /bin/nc * -e /bin/bash* > cron && crontab cron*'
    - '* /bin/nc * -e /bin/bash*> * crontab cron*'
  - process.command_line:
    - '*nc -l -p * -e cmd.exe*'
    - '*nc -l -p * -e *.bat*'
    - '*nc -l -p * -e /bin/bash*'
    - '*echo nc -l -p * > *.bat*'
    - '*nc -v -n -z -w1 *-*'
    - '*nc -z -v * *'
    - '* /bin/nc * -e /bin/bash* > cron && crontab cron*'
    - '* /bin/nc * -e /bin/bash*> * crontab cron*'
  - process.env_vars:
    - '*nc -l -p * -e cmd.exe*'
    - '*nc -l -p * -e *.bat*'
    - '*nc -l -p * -e /bin/bash*'
    - '*echo nc -l -p * > *.bat*'
    - '*nc -v -n -z -w1 *-*'
    - '*nc -z -v * *'
    - '* /bin/nc * -e /bin/bash* > cron && crontab cron*'
    - '* /bin/nc * -e /bin/bash*> * crontab cron*'
  - process.io.text:
    - '*nc -l -p * -e cmd.exe*'
    - '*nc -l -p * -e *.bat*'
    - '*nc -l -p * -e /bin/bash*'
    - '*echo nc -l -p * > *.bat*'
    - '*nc -v -n -z -w1 *-*'
    - '*nc -z -v * *'
    - '* /bin/nc * -e /bin/bash* > cron && crontab cron*'
    - '* /bin/nc * -e /bin/bash*> * crontab cron*'
  - process.title:
    - '*nc -l -p * -e cmd.exe*'
    - '*nc -l -p * -e *.bat*'
    - '*nc -l -p * -e /bin/bash*'
    - '*echo nc -l -p * > *.bat*'
    - '*nc -v -n -z -w1 *-*'
    - '*nc -z -v * *'
    - '* /bin/nc * -e /bin/bash* > cron && crontab cron*'
    - '* /bin/nc * -e /bin/bash*> * crontab cron*'
  - pe.company:
    - '*nc -l -p * -e cmd.exe*'
    - '*nc -l -p * -e *.bat*'
    - '*nc -l -p * -e /bin/bash*'
    - '*echo nc -l -p * > *.bat*'
    - '*nc -v -n -z -w1 *-*'
    - '*nc -z -v * *'
    - '* /bin/nc * -e /bin/bash* > cron && crontab cron*'
    - '* /bin/nc * -e /bin/bash*> * crontab cron*'
  - pe.description:
    - '*nc -l -p * -e cmd.exe*'
    - '*nc -l -p * -e *.bat*'
    - '*nc -l -p * -e /bin/bash*'
    - '*echo nc -l -p * > *.bat*'
    - '*nc -v -n -z -w1 *-*'
    - '*nc -z -v * *'
    - '* /bin/nc * -e /bin/bash* > cron && crontab cron*'
    - '* /bin/nc * -e /bin/bash*> * crontab cron*'
  - pe.product:
    - '*nc -l -p * -e cmd.exe*'
    - '*nc -l -p * -e *.bat*'
    - '*nc -l -p * -e /bin/bash*'
    - '*echo nc -l -p * > *.bat*'
    - '*nc -v -n -z -w1 *-*'
    - '*nc -z -v * *'
    - '* /bin/nc * -e /bin/bash* > cron && crontab cron*'
    - '* /bin/nc * -e /bin/bash*> * crontab cron*'
  - os.full:
    - '*nc -l -p * -e cmd.exe*'
    - '*nc -l -p * -e *.bat*'
    - '*nc -l -p * -e /bin/bash*'
    - '*echo nc -l -p * > *.bat*'
    - '*nc -v -n -z -w1 *-*'
    - '*nc -z -v * *'
    - '* /bin/nc * -e /bin/bash* > cron && crontab cron*'
    - '* /bin/nc * -e /bin/bash*> * crontab cron*'
  - host.hostname:
    - '*nc -l -p * -e cmd.exe*'
    - '*nc -l -p * -e *.bat*'
    - '*nc -l -p * -e /bin/bash*'
    - '*echo nc -l -p * > *.bat*'
    - '*nc -v -n -z -w1 *-*'
    - '*nc -z -v * *'
    - '* /bin/nc * -e /bin/bash* > cron && crontab cron*'
    - '* /bin/nc * -e /bin/bash*> * crontab cron*'
  - event.original:
    - '*nc -l -p * -e cmd.exe*'
    - '*nc -l -p * -e *.bat*'
    - '*nc -l -p * -e /bin/bash*'
    - '*echo nc -l -p * > *.bat*'
    - '*nc -v -n -z -w1 *-*'
    - '*nc -z -v * *'
    - '* /bin/nc * -e /bin/bash* > cron && crontab cron*'
    - '* /bin/nc * -e /bin/bash*> * crontab cron*'
  - email.subject:
    - '*nc -l -p * -e cmd.exe*'
    - '*nc -l -p * -e *.bat*'
    - '*nc -l -p * -e /bin/bash*'
    - '*echo nc -l -p * > *.bat*'
    - '*nc -v -n -z -w1 *-*'
    - '*nc -z -v * *'
    - '* /bin/nc * -e /bin/bash* > cron && crontab cron*'
    - '* /bin/nc * -e /bin/bash*> * crontab cron*'
  - device.model.name:
    - '*nc -l -p * -e cmd.exe*'
    - '*nc -l -p * -e *.bat*'
    - '*nc -l -p * -e /bin/bash*'
    - '*echo nc -l -p * > *.bat*'
    - '*nc -v -n -z -w1 *-*'
    - '*nc -z -v * *'
    - '* /bin/nc * -e /bin/bash* > cron && crontab cron*'
    - '* /bin/nc * -e /bin/bash*> * crontab cron*'
  - container.image.name:
    - '*nc -l -p * -e cmd.exe*'
    - '*nc -l -p * -e *.bat*'
    - '*nc -l -p * -e /bin/bash*'
    - '*echo nc -l -p * > *.bat*'
    - '*nc -v -n -z -w1 *-*'
    - '*nc -z -v * *'
    - '* /bin/nc * -e /bin/bash* > cron && crontab cron*'
    - '* /bin/nc * -e /bin/bash*> * crontab cron*'
  - container.name:
    - '*nc -l -p * -e cmd.exe*'
    - '*nc -l -p * -e *.bat*'
    - '*nc -l -p * -e /bin/bash*'
    - '*echo nc -l -p * > *.bat*'
    - '*nc -v -n -z -w1 *-*'
    - '*nc -z -v * *'
    - '* /bin/nc * -e /bin/bash* > cron && crontab cron*'
    - '* /bin/nc * -e /bin/bash*> * crontab cron*'
  condition: selection
fields:
- ParentCommandLine
- CommandLine
- Details
- registry.value
- process.args
- process.command_line
- process.env_vars
- process.io.text
- process.title
- pe.company
- pe.description
- pe.product
- os.full
- host.hostname
- event.original
- email.subject
- device.model.name
- container.image.name
- container.name
falsepositives:
- unknown
level: medium
